# libexec/rc/etc.init.d/bsnmpd
#!/sbin/openrc-run
#
# $FreeBSD$
#

depend()
{
	need syslogd
	keyword -jail -shutdown
}

name="bsnmpd"
desc="Simple and extensible SNMP daemon"
command="/usr/sbin/${name}"
pidfile="${bsnmpd_pidfile:-/var/run/snmpd.pid}"
command_args="-p ${pidfile}"
