# libexec/rc/etc.init.d/devd
#!/sbin/openrc-run

command=/sbin/devd
command_args="$devd_args -f /etc/devd-openrc.conf"
name="devd"
description="Device State Change Daemon"

depend() {
	after network ldconfig
	keyword -nojail -shutdown
}
