# libexec/rc/etc.init.d/gssd
#!/sbin/openrc-run

command=/usr/sbin/gssd
name="Generic Security Services Daemon"

depend()
{
	need localmount
}
