# libexec/rc/etc.init.d/hostapd
#!/sbin/openrc-run

name="hostapd"
description="Authenticator for IEEE 802.11 networks"
command=/usr/sbin/hostapd
ifn="$2"
if [ -z "$ifn" ]; then
	rcvar="hostapd_enable"
	conf_file="/etc/${name}.conf"
	pidfile="/var/run/${name}.pid"
else
	rcvar=
	conf_file="/etc/${name}-${ifn}.conf"
	pidfile="/var/run/${name}-${ifn}.pid"
fi
command_args="-P ${pidfile} -B ${conf_file}"
required_files="${conf_file}"


depend()
{
	need localmount
	use net
	keyword -shutdown -jail
}

start_pre()
{
	load_kld wlan_xauth wlan_wep wlan_tkip wlan_ccmp
}
