# libexec/rc/etc.init.d/iscsid
#!/sbin/openrc-run

depend()
{
        use net
}

name="iscsid"
desc="iSCSI initiator daemon"
pidfile="/var/run/${name}.pid"
command="/usr/sbin/${name}"
required_modules="iscsi"
