# libexec/rc/etc.init.d/keyserv
#!/sbin/openrc-run

name="keyserv"
command="/usr/sbin/${name}"

depend() {
	need ypset rpcbind
	keyword -shutdown
}
