# libexec/rc/etc.init.d/ldapd
#!/sbin/openrc-run

name="ldapd"
rcvar="ldapd_enable"
command="/usr/sbin/${name}"

depend() {
	need devfs
	keyword -jail -shutdown
}
