# libexec/rc/etc.init.d/lockd
#!/sbin/openrc-run

command=/usr/sbin/rpc.lockd
command_args=$lockd_args
name="lockd"

depend()
{
	need localmount rpcbind statd
	use net logger dns
	before inetd xinetd
}
