# libexec/rc/etc.init.d/mdconfig
#!/sbin/openrc-run
# OpenRC -> rc.d redirect script
# Created by Ken Moore <ken@ixsystems.com>
name="mdconfig"
desc="Create and control memory disks"

rcfile="/etc/rc.d/${name}"

depend() {
	before abi ldconfig savecore watchdogd
	need swap root
}

start() {
  if [ $(sysrc -n "${name}_enable") = "YES" ] ; then
    ${rcfile} start
  else
    ${rcfile} onestart
  fi
  return 0
}

stop() {
  if [ $(sysrc -n "${name}_enable") = "YES" ] ; then
    ${rcfile} stop
  else
    ${rcfile} onestop
  fi
  return 0
}
