# libexec/rc/etc.init.d/nfscbd
#!/sbin/openrc-run
#

name="nfscbd"
command="/usr/sbin/${name}"

depend() {
	need nfsuserd
	keyword -jail -shutdown
}
