# libexec/rc/etc.init.d/rwho
#!/sbin/openrc-run

name="rwhod"
rcvar="rwhod_enable"
command="/usr/sbin/${name}"

depend() {
	before login
	keyword -shutdown
}
