# libexec/rc/etc.init.d/timed
#!/sbin/openrc-run

name="timed"
desc="Time server daemon"
command="/usr/sbin/${name}"

depend() {
	before login
	keyword -jail -shutdown
}
