# libexec/rc/etc.init.d/watchdogd
#!/sbin/openrc-run

name="watchdogd"
command="/usr/sbin/${name}"
pidfile="/var/run/${name}.pid"

depend() {
	after syslogd
	keyword -jail -shutdown
}
