# libexec/rc/etc.init.d/zfs
#!/sbin/openrc-run

name="zfs"
desc="Dummy service for zfs services"

depend()
{
	need zfsbe zfsd zfsmount zpool zvol
	before bootmisc logger network
	after clock root sysctl
	keyword -prefix
}

start()
{
	return 0
}