# libexec/rc/etc.init.d/zfsd
#!/sbin/openrc-run

name="zfsd"
command="/usr/sbin/${name}"

depend() {
	after devd zfs
	keyword -shutdown
}
