# libexec/rc/etc.init.d/zfsmount
#!/sbin/openrc-run

depend()
{
	after zfsbe
	before logger network
}

. /etc/rc.subr

PATH="${PATH}:/usr/local/sbin:/usr/local/bin"
export PATHH

checksystem() {
	if [ ! -c /dev/zfs ]; then
		if [ ! -x /boot/kernel/zfs.ko -a ! -x /boot/modules/zfs.ko ]; then
			eerror "zfs module not found. Is the kernel built with ZFS support"
			return 1
		fi
	fi
	which -s zfs
	if [ $? -ne 0 ]; then
		eerror "zfs binary not found."
		return 1
	fi
	return 0
}


zfs_start_jail()
{
	if [ `$SYSCTL_N security.jail.mount_allowed` -eq 1 ]; then
		zfs mount -a
	fi
}

zfs_start_main()
{
	zfs mount -va
	zfs share -a
	if [ ! -r /etc/zfs/exports ]; then
		touch /etc/zfs/exports
	fi
}

start()
{
	ebegin "Starting ZFS"
	checksystem || return 1

	if [ `$SYSCTL_N security.jail.jailed` -eq 1 ]; then
		zfs_start_jail
	else
		zfs_start_main
	fi
	eend 0
}

zfs_stop_jail()
{
	if [ `$SYSCTL_N security.jail.mount_allowed` -eq 1 ]; then
		zfs unmount -a
	fi
}

zfs_stop_main()
{
	zfs unshare -a
	zfs unmount -a
}

stop()
{
	ebegin "Unsharing ZFS filesystems"
	if [ `$SYSCTL_N security.jail.jailed` -eq 1 ]; then
		zfs_stop_jail
	else
		zfs_stop_main
	fi
	eend 0
}
