# libexec/rc/etc.init.d/zpool
#!/sbin/openrc-run

depend()
{
	use hostid
	before bootmisc
	keyword -jail
}

. /etc/rc.subr

PATH="${PATH}:/usr/local/sbin:/usr/local/bin"
export PATH

checksystem() {
	if [ ! -c /dev/zfs ]; then
		if [ ! -x /boot/kernel/zfs.ko -a ! -x /boot/modules/zfs.ko ]; then
			eerror "zfs module not found. Is the kernel built with ZFS support"
			return 1
		fi
	fi
	which -s zpool
	if [ $? -ne 0 ]; then
		 eerror "zpool binary not found."
		 return 1
    fi
	return 0
}

start()
{
	ebegin "Starting zpool"
	checksystem || return 1

	local cachefile
	for cachefile in /etc/zfs/zpool.cache /boot/zfs/zpool.cache; do
		if [ -r $cachefile ]; then
			zpool import -c $cachefile -a -N
			if [ $? -ne 0 ]; then
				echo "Import of zpool cache ${cachefile} failed," \
				    "will retry after root mount hold release"
				root_hold_wait
				zpool import -c $cachefile -a -N
			fi
			break
		fi
	done
	eend 0
}


status()
{
	# show pool status and list
	zpool status && echo && zpool list
}