# libexec/rc/etc.init.d/zvol
#!/sbin/openrc-run

depend()
{
	after zpool
	keyword -jail -prefix
}

. /etc/rc.subr

PATH="${PATH}:/usr/local/sbin:/usr/local/bin"
export PATH

ZFS="zfs"
ZPOOL="zpool"
ZPOOL_CACHE="/boot/zfs/zpool.cache"
ZFS_MODULE=zfs

checksystem() {
	if [ ! -c /dev/zfs ]; then
		if [ ! -x /boot/kernel/zfs.ko -a ! -x /boot/modules/zfs.ko ]; then
			eerror "zfs module not found. Is the kernel built with ZFS support"
			return 1
		fi
	fi
	which -s zfs
	if [ $? -ne 0 ]; then
		eerror "zfs binary not found."
		return 1
	fi
	return 0
}

start() {
	ebegin "Starting swap on ZVOLs"
	checksystem || return 1

	# Enable swap on ZVOLs with property org.freebsd:swap=on.
	zfs list -H -o org.freebsd:swap,name -t volume | \
	while read state name; do
		case "${state}" in
		[oO][nN])
			swapon /dev/zvol/${name}
			;;
		esac
	done
	eend 0
}

stop()
{
	ebegin "Stopping swap on ZVOLs"
	zfs list -H -o org.freebsd:swap,name -t volume | \
        while read state name; do
                case "${state}" in
                [oO][nN])
                        swapoff /dev/zvol/${name}
                        ;;
                esac
        done
	eend 0
}
