# libexec/rc/etc.init.d/zvol
#!/sbin/openrc-run
# Copyright 1999-2011 Gentoo Foundation
# Released under the 2-clause BSD license.
# $Header: /var/cvsroot/gentoo-x86/sys-fs/zfs/files/zfs,v 0.9 2011/04/30 10:13:43 devsk Exp $

depend()
{
	after zfs
	keyword -jail -prefix
}

PATH="${PATH}:/usr/local/bin:/usr/local/sbin"
export PATH

ZFS="zfs"

start() {
	ebegin "Starting ZFS swap"
	# Enable swap on ZVOLs with property org.freebsd:swap=on.
	${ZFS} list -H -o org.freebsd:swap,name -t volume | \
	while read state name; do
		case "${state}" in
		[oO][nN])
			swapon /dev/zvol/${name}
			;;
		esac
	done
	eend 0
}

stop()
{	
	ebegin "Stopping ZFS swap"
	${ZFS} list -H -o org.freebsd:swap,name -t volume | \
        while read state name; do
                case "${state}" in
                [oO][nN])
                        swapoff /dev/zvol/${name}
                        ;;
                esac
        done
	eend 0
}
